# systemd/pdf_audio_kitoblar_bot-cleanup.service
[Unit]
Description=pdf_audio_kitoblar_bot - Cleanup Worker
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
# Template unit for public repo use.
# Replace the generic user/group/path values to match your host before installing.
User=pdfbot
Group=pdfbot
WorkingDirectory=/opt/pdf_audio_kitoblar_bot
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/etc/pdf_audio_kitoblar_bot.env
ExecStart=/opt/pdf_audio_kitoblar_bot/.venv/bin/python -m workers.cleanup_worker
Restart=always
RestartSec=10
TimeoutStopSec=20
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
