# systemd/pdf_audio_kitoblar_bot-dashboard.service
[Unit]
Description=pdf_audio_kitoblar_bot - Local Admin Dashboard
Wants=network-online.target pdf_audio_kitoblar_bot.service pdf_audio_kitoblar_bot-bot.service
After=network-online.target pdf_audio_kitoblar_bot.service pdf_audio_kitoblar_bot-bot.service
PartOf=pdf_audio_kitoblar_bot-stack.target pdf_audio_kitoblar_bot-bot.service

[Service]
Type=simple
# Template unit for public repo use.
# Replace the generic user/group/path values to match your host before installing.
User=pdfbot
Group=pdfbot
WorkingDirectory=/opt/pdf_audio_kitoblar_bot
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/etc/pdf_audio_kitoblar_bot.env
ExecStart=/opt/pdf_audio_kitoblar_bot/.venv/bin/python /opt/pdf_audio_kitoblar_bot/dashboard_server.py --host 127.0.0.1 --port 8091
Restart=always
RestartSec=3
TimeoutStopSec=20
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
